# croma-wifi.service
[Unit]
Description=Croma WiFi Manager - Hotspot Automático
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/Chromasistem
ExecStart=/usr/bin/python3 /home/pi/Chromasistem/wifi_manager.py monitor
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
